# deploy/ansible/roles/honeyagent/templates/honeyagent.service.j2
[Unit]
Description=Honeyagent deception service
After=network.target

[Service]
Type=simple
User=honeyagent
Group=honeyagent
ExecStart=/usr/local/bin/honeyagent --config /etc/honeyagent/config.json
Restart=on-failure
RestartSec=5

AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/log/honeyagent /var/log/kaspersky /var/opt/kaspersky

[Install]
WantedBy=multi-user.target
